# doc/rsync_backup.crontab.example
0 0 * * * bash -c 'sleep $(($RANDOM\%120))m' && ~/opt/rsyncbak/bin/rsync_backup
0 0 * * * bash -c 'sleep $(($RANDOM\%120))m' && ~/opt/rsyncbak/bin/mysql_backup
